# config/claude-mobile.service
[Unit]
Description=Claude Mobile Web Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/claude-mobile/server
ExecStart=/usr/bin/node /home/%i/claude-mobile/server/index.js
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

Environment=NODE_ENV=production
Environment=PORT=3000

# Logging
StandardOutput=journal
StandardError=journal

# Security hardening
# NOTE: NoNewPrivileges is intentionally absent — tmux child processes need to
# be able to run sudo commands inside the terminal sessions.
PrivateTmp=true
ProtectSystem=strict
# node-pty needs /dev/ptmx so PrivateDevices must NOT be set.
# All user-writable paths (repos, config, sessions, app) are under /home/%i.
ReadWritePaths=/home/%i

[Install]
WantedBy=multi-user.target
