# config/claude-mobile.service
[Unit]
Description=Claude Mobile Web Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/claude-mobile/server

# ─── Node.js flags (generated by setup.sh — profile: e2-micro) ───────
ExecStart=/usr/bin/node --max-old-space-size=256 --optimize-for-size --expose-gc --max-semi-space-size=2 /home/%i/claude-mobile/server/index.js

Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

Environment=NODE_ENV=production
Environment=PORT=3000
# UV_THREADPOOL_SIZE: ridotto su VM con poca RAM, standard altrove
Environment=UV_THREADPOOL_SIZE=2

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-mobile

# ─── cgroup resource limits (generated by setup.sh — profile: e2-micro) ─
# MemoryHigh: soft limit — kernel throttles above this
# MemoryMax:  hard limit — OOM kill (prefer controlled restart via Restart=always)
MemoryAccounting=true
MemoryHigh=400M
MemoryMax=512M

# Prevent fork bombs or runaway PTY spawning
TasksMax=200

# Prefer killing other processes when OOM killer activates
OOMScoreAdjust=-500

# I/O: slightly boost over default for responsive terminal
IOSchedulingClass=best-effort
IOSchedulingPriority=4

# ─── Security hardening ──────────────────────────────────────────────────────
# NOTE: NoNewPrivileges is intentionally absent — tmux sessions may need sudo.
PrivateTmp=true
ProtectSystem=strict
# node-pty needs /dev/ptmx so PrivateDevices must NOT be set.
ReadWritePaths=/home/%i

# Limit open file descriptors
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
